# supervisord.conf
[supervisord]
nodaemon=true

[program:openvpn]
command=/usr/local/bin/openvpn.sh
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
